# jobs/jobs.cron
00 12 * * * php /var/www/html/PoliticallyMail/jobs/updateDictionary.php
